# deploy/indus-ure-backup.service
[Unit]
Description=Daily PostgreSQL backup for INDUS URE
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
User=postgres
Group=postgres
ExecStart=/usr/local/sbin/indus-ure-backup
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/backups/indus-ure
