# deploy/indus-ure.service
[Unit]
Description=INDUS URE Node.js application
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=indus-ure
Group=indus-ure
WorkingDirectory=/opt/indus-ure/current
Environment=NODE_ENV=production
EnvironmentFile=/etc/indus-ure.env
ExecStart=/usr/bin/node /opt/indus-ure/current/outputs/server.js
Restart=on-failure
RestartSec=5s
TimeoutStopSec=20s
UMask=0077

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadWritePaths=/var/lib/indus-ure
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target